# docker-mt5-node/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor

[program:xvfb]
command=Xvfb :0 -screen 0 1280x800x24
priority=10
autorestart=true
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb-err.log

[program:fluxbox]
command=fluxbox -display :0
priority=20
autorestart=true
environment=DISPLAY=":0"
stdout_logfile=/var/log/supervisor/fluxbox.log
stderr_logfile=/var/log/supervisor/fluxbox-err.log

[program:x11vnc]
command=x11vnc -display :0 -nopw -forever -shared -rfbport 5900
priority=30
autorestart=true
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc-err.log

[program:novnc]
command=websockify --web=/usr/share/novnc/ 8080 localhost:5900
priority=40
autorestart=true
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc-err.log

[program:mt5-install]
command=bash -c "MT5EXE='/home/wineuser/.wine/drive_c/Program Files/MetaTrader 5/terminal64.exe'; if [ ! -f \"$MT5EXE\" ]; then echo '[MT5] Installing MetaTrader 5...'; DISPLAY=:0 wine64 /home/wineuser/mt5setup.exe /auto; echo '[MT5] Install finished.'; else echo '[MT5] Already installed.'; fi"
priority=50
autorestart=false
startsecs=0
stdout_logfile=/var/log/supervisor/mt5-install.log
stderr_logfile=/var/log/supervisor/mt5-install-err.log
environment=DISPLAY=":0",WINEPREFIX="/home/wineuser/.wine"

[program:mt5]
command=bash -c "MT5EXE='/home/wineuser/.wine/drive_c/Program Files/MetaTrader 5/terminal64.exe'; echo '[MT5] Waiting for terminal64.exe...'; for i in $(seq 1 120); do [ -f \"$MT5EXE\" ] && break; echo \"[MT5] Waiting... ($i/120)\"; sleep 5; done; if [ -f \"$MT5EXE\" ]; then echo '[MT5] Starting MetaTrader 5...'; wine64 \"$MT5EXE\" /portable; else echo '[MT5] ERROR: terminal64.exe not found after 10 min.'; exit 1; fi"
priority=60
autorestart=true
startsecs=10
startretries=5
stdout_logfile=/var/log/supervisor/mt5.log
stderr_logfile=/var/log/supervisor/mt5-err.log
environment=DISPLAY=":0",WINEPREFIX="/home/wineuser/.wine"
