# docker-mt5-node/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=Xvfb :0 -screen 0 1280x800x24
priority=10
autorestart=true

[program:fluxbox]
command=fluxbox -display :0
priority=20
autorestart=true
environment=DISPLAY=":0"

[program:x11vnc]
command=x11vnc -display :0 -nopw -forever -shared -bg
priority=30
autorestart=true

[program:novnc]
command=websockify --web=/usr/share/novnc/ 8080 localhost:5900
priority=40
autorestart=true

[program:mt5-install]
# Only runs once, if MT5 isn't installed. Otherwise it just exits.
command=bash -c "if [ ! -d '/root/.wine/drive_c/Program Files/MetaTrader 5' ]; then wine /root/mt5setup.exe /auto; fi"
priority=50
autorestart=false
environment=DISPLAY=":0",WINEPREFIX="/root/.wine"

[program:mt5]
command=bash -c "sleep 10 && wine '/root/.wine/drive_c/Program Files/MetaTrader 5/terminal64.exe' /portable"
priority=60
autorestart=true
environment=DISPLAY=":0",WINEPREFIX="/root/.wine"
